# spondkassier.service
[Unit]
Description=Spond Kassier Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=schoasda
WorkingDirectory=/home/schoasda/SpondKassierBot
ExecStart=/home/schoasda/SpondKassierBot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
